# translation-service.service
[Unit]
Description=Automated translation service for Java property files
After=network.target

[Service]
Type=oneshot
User=bisquser
WorkingDirectory=/home/bisquser/workspace/translate-java-property-files
ExecStart=/bin/bash /home/bisquser/workspace/translate-java-property-files/deploy.sh
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/home/bisquser/.local/bin"
EnvironmentFile=/home/bisquser/.env
# Add GitHub token for PR creation
Environment="GITHUB_TOKEN=your_token_will_be_set_in_environment"

[Install]
WantedBy=multi-user.target 